# root/etc/services.d/deluge-web/run
#!/usr/bin/with-contenv bash

umask 022

cmd="s6-setuidgid abc /usr/bin/deluge-web -c ${CONFIG} --loglevel=info"
echo "Executing: [$cmd]"
exec $cmd
