# root/etc/services.d/deluged/run
#!/usr/bin/with-contenv bash

umask 022

exec \
	s6-setuidgid abc /usr/bin/deluged -c "${CONFIG}" \
	-d --loglevel=info -l "${CONFIG}/deluged.log"
